# init/gateway.conf
[supervisord]
nodaemon=true

[program:gateway]
command=./gateway

autostart=true
autorestart=true
startretries=10

stdout_logfile=gateway.log
stdout_logfile_maxbytes=0