# init/gateway.conf
[supervisord]
nodaemon=true

[program:gateway]
command=/usr/local/bin/gateway

autostart=true
autorestart=true

stdout_logfile=stdout.log
stderr_logfile=stderr.log