# init/gmicro.conf
[supervisord]
nodaemon=true

[program:gmicro]
command=/usr/local/bin/gmicro

autostart=true
autorestart=true

stdout_logfile=stdout.log
stderr_logfile=stderr.log