# init/tmicro.conf
[supervisord]
nodaemon=true

[program:gmicro]
directory=/usr/local
command=/usr/local/bin/tmicro

autostart=true
autorestart=true